# deploy/systemd/noogh-api.service
[Unit]
Description=Noogh Unified System - API Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=noogh
Group=noogh
WorkingDirectory=/home/noogh/projects/noogh_unified_system
Environment="PATH=/home/noogh/projects/noogh_unified_system/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/noogh/projects/noogh_unified_system"
Environment="API_HOST=0.0.0.0"
Environment="API_PORT=8000"
Environment="LOG_LEVEL=INFO"

ExecStart=/home/noogh/projects/noogh_unified_system/venv/bin/uvicorn src.api.main:app --host 0.0.0.0 --port 8000 --workers 2

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=append:/home/noogh/projects/noogh_unified_system/logs/api.log
StandardError=append:/home/noogh/projects/noogh_unified_system/logs/api-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
